# systemd/afterburn-sshkeys@.service.in
[Unit]
Description=Afterburn (SSH Keys)
# These platforms support SSH keys and provide a consistent source for them.
# There may be additional platforms which support SSH keys only in some cases
# (e.g. via optional platform components); those platforms need a user-provided
# dropin, adding an appropriate triggering condition and setting the value of
# `AFTERBURN_OPT_PROVIDER` as needed.
ConditionKernelCommandLine=|ignition.platform.id=aliyun
ConditionKernelCommandLine=|ignition.platform.id=aws
ConditionKernelCommandLine=|ignition.platform.id=azure
ConditionKernelCommandLine=|ignition.platform.id=digitalocean
ConditionKernelCommandLine=|ignition.platform.id=exoscale
ConditionKernelCommandLine=|ignition.platform.id=gcp
ConditionKernelCommandLine=|ignition.platform.id=one
ConditionKernelCommandLine=|ignition.platform.id=packet
ConditionKernelCommandLine=|ignition.platform.id=vultr
ConditionKernelCommandLine=|ignition.platform.id=ibmcloud

[Service]
Type=oneshot
Environment=AFTERBURN_OPT_PROVIDER=--cmdline
ExecStart=/usr/bin/afterburn ${AFTERBURN_OPT_PROVIDER} --ssh-keys=%i

[Install]
DefaultInstance=@DEFAULT_INSTANCE@
RequiredBy=afterburn-sshkeys.target
