# systemd/afterburn-sshkeys.target
[Unit]
Description=Synchronize afterburn-sshkeys@.service template instances

[Install]
RequiredBy=multi-user.target
